# deployment/files/conf/gunicorn_nessmap.conf.j2
[program:nessmap]
command = /var/projects/nessmap/env/bin/gunicorn_django -b 0.0.0.0:8080
directory = /var/projects/nessmap/app/nessmap
user = nessmap